# templates/iqbus-udev-handler.service.template
[Unit]
Description=Handle iqbus RTL-SDR udev event
After=systemd-udevd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/iqbus-udev-handler.sh
